# packaging/systemd/tardigrade.service
[Unit]
Description=Tardigrade edge gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tardigrade
Group=tardigrade
WorkingDirectory=/var/lib/tardigrade
EnvironmentFile=-/etc/tardigrade/tardigrade.env
ExecStart=/usr/bin/tardi run -c /etc/tardigrade/tardigrade.conf
ExecReload=/usr/bin/tardi reload -c /etc/tardigrade/tardigrade.conf
ExecStop=/usr/bin/tardi stop -c /etc/tardigrade/tardigrade.conf
Restart=on-failure
RestartSec=2
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
